# res/daf.conf
[program:daf_proccess]
command=/home/pi/daf-pi/venv/bin/python3 /home/pi/daf-pi/app.py
directory=/home/pi/daf-pi
stdout_logfile=/home/pi/daf-log.txt
redirect_stderr=true